# .crontab
12 10 * * 1-5   open https://global.gotomeeting.com/join/629430197
